# profile.d/java.sh

export JAVA_HOME="/usr/lib/jvm/java-1.8.0-openjdk/bin/java"
export JAVA="/usr/lib/jvm/java-1.8.0-openjdk/bin/java"
